# talkat.service
[Unit]
Description=Talkat Voice Dictation Server
Documentation=https://github.com/ronakrm/talkat
After=sound.target

[Service]
Type=simple
ExecStart=/usr/bin/talkat server
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening — defence-in-depth for a per-user service. Kept in sync
# with the unit `talkat install-service` writes (src/talkat/service.py).
# ReadWritePaths punches holes through ProtectHome=read-only for the
# talkat-owned XDG directories. XDG_RUNTIME_DIR (/run/user/$UID) is
# outside $HOME so it stays writable without an explicit entry.
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.cache/talkat %h/.local/share/talkat %h/.config/talkat
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
# PrivateDevices stays off — /dev/snd access matters if a future server-side
# capture path lands; today only the client touches the mic.
PrivateDevices=no

[Install]
WantedBy=default.target
