# talkat.service
[Unit]
Description=Talkat Voice Dictation Server
Documentation=https://github.com/yourusername/talkat
After=pulseaudio.service pipewire.service
Wants=pulseaudio.service

[Service]
Type=simple
ExecStart=%i/bin/talkat server
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=HOME=%h
Environment=XDG_CONFIG_HOME=%h/.config
Environment=XDG_CACHE_HOME=%h/.cache
Environment=XDG_DATA_HOME=%h/.local/share

# Security hardening for user service
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=%h/.cache/talkat %h/.local/share/talkat %h/.config/talkat
NoNewPrivileges=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
PrivateDevices=no  # Need access to audio devices

[Install]
WantedBy=default.target